# backend/infrastructure/cron/backend-cron@update-expire-at.timer.d/override.conf
[Timer]
# Clear the template schedule
OnCalendar=

# Custom schedule for update-expire-at
OnCalendar=*-*-* *:00:00
