# backend/infrastructure/cron/backend-cron@sync-completed.timer.d/override.conf
[Timer]
# Clear the template schedule
OnCalendar=

# Custom schedule for sync-completed
OnCalendar=*-*-* *:0/2:00
