# backend/infrastructure/cron/staging/backend-staging-cron@check-password-expiration.timer.d/override.conf
[Timer]
# Clear the template schedule
OnCalendar=

# Custom schedule for check-password-expiration
OnCalendar=*-*-* *:0/15:00
