# backend/infrastructure/cron/backend-cron@update-rank.timer.d/override.conf
[Timer]
# Clear the template schedule
OnCalendar=

# Custom schedule for update-rank
OnCalendar=*-*-* 00:00:00
