# backend/infrastructure/cron/staging/backend-staging-cron@sync-account-availability.timer.d/override.conf
[Timer]
# Clear the template schedule
OnCalendar=

# Custom schedule for sync-account-availability
OnCalendar=*-*-* *:0/2:00
